# examples/upstart/neckup.conf
description "neckup web service"
author "William Wennerström <william@wiol.io>"

start on (net-device-up
	  and local-filesystems
	  and runlevel [2345])

stop on runlevel [016]

respawn

exec start-stop-daemon --start -d /home/neckup/neckup --chuid neckup --exec ./neckup
